# deploy/mtg-agent.service
[Unit]
Description=MTG Agent MCP Server
After=network.target mongod.service
Requires=mongod.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/mtg_agent
EnvironmentFile=/home/admin/mtg_agent/.env
ExecStart=/home/admin/mtg_agent/.venv/bin/mtg-agent
Restart=on-failure
RestartSec=5

# stdout/stderr go to journald — view with: journalctl -u mtg-agent -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
